# watchdogs/disk-temp-watchdog/disk-temp-watchdog.service
[Unit]
Description=CNJ Mesh disk and temperature watchdog
After=network.target

[Service]
Type=oneshot
Environment=CNJ_DISCORD_WEBHOOK=REPLACE_ME
Environment=NODE_LABEL=REPLACE_ME_NODE_LABEL
ExecStart=/usr/bin/python3 /opt/disk-temp-watchdog/watchdog.py
